# sshd-core/src/test/resources/org/apache/sshd/client/opensshcerts/docker/supervisord.conf
[supervisord]
; nodaemon=true (no daemon) means run in foreground
nodaemon=true
user=root

[program:sshd]
directory=/etc/ssh
command=/usr/sbin/sshd -D -e
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0